# supervisord.conf
[supervisord]
nodaemon=true
user=fireuser
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:iot-service]
command=python Main.py
directory=/app
user=fireuser
autostart=%(ENV_START_IOT)s
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/iot-service.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1

[program:web-service]
command=gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 App:app
directory=/app
user=fireuser
autostart=%(ENV_START_WEB)s
autorestart=true
redirect_stderr=true
stdout_logfile=/tmp/web-service.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED=1,FLASK_ENV=production

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface